# php/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/local/bin/apache2-foreground

[program:crond]
command=/bin/busybox crond -f -L /dev/stderr
